# autostart/randombg.service
[Unit]
Description=RandomBG Tray App
After=network.target

[Service]
Type=simple
# Passe diesen Pfad an deinen Projekt- oder Installationsort an
WorkingDirectory=%h/RandomBG
# Verwende dein (virtuelles) Environment, falls vorhanden
Environment=PATH=%h/RandomBG/.venv/bin:/usr/local/bin:/usr/bin
# Übernimmt DISPLAY, WAYLAND_DISPLAY und XAUTHORITY aus deiner Session.
# Stelle sicher, dass der User-Service mit einer aktiven Session gestartet
# wird oder importiere die Variablen einmalig mit:
#   systemctl --user import-environment DISPLAY WAYLAND_DISPLAY XAUTHORITY
ExecStart=/usr/bin/env python -m random_bg.app
Restart=on-failure

[Install]
WantedBy=default.target
